# .config/systemd/user/wallpaper-changer.service
[Unit]
Description="Change the wallpaper"

[Service]
Environment="DISPLAY=:0"
ExecStart=%h/.local/bin/_wallpaper_changer.sh
